# deploy/systemd/codex-multi-account.service
# 这个文件用于把 codex-multi-account 作为本机 systemd 服务常驻运行。

[Unit]
Description=codex-multi-account local service
After=network.target

[Service]
Type=simple
User=djy
WorkingDirectory=/home/djy/codex-multi-account
Environment=CMA_DATA_DIR=/home/djy/codex-multi-account/data
ExecStart=/bin/bash /home/djy/codex-multi-account/scripts/run-codex-multi-account.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
